# awdd.service
[Service]
ExecStart=/usr/local/sbin/awdd.sh

[Install]
WantedBy=default.target
